# setup_tmux/systemd/user/claude-auto-approve.service
[Unit]
Description=Auto-approve safe Claude Code permission prompts
After=default.target

[Service]
ExecStart=%h/.local/bin/claude-auto-approve
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
