# setup_tmux/systemd/user/tmux-compose-snapshot.service
[Unit]
Description=Snapshot tmux session layout for tmux-compose recovery

[Service]
Type=oneshot
ExecStart=%h/.local/bin/tmux-compose snapshot
